# etc/s6/dnsmasq/run
#!/bin/sh

exec /usr/sbin/dnsmasq -C "$DNSMASQ_CONF"
