# etc/s6/tor/run
#!/bin/sh

exec /usr/sbin/s6-setuidgid tor tor --defaults-torrc "$TOR_CONF"
